# mypage/app/docker/supervisor-app.conf
[program:app-uwsgi]
command = /root/.pyenv/shims/uwsgi --ini /home/docker/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
